# systemd/watchtower-oneshot.service
[Unit]
Description=Watchtower one-time update run
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/run-watchtower.sh
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
